# cmd/swiftwave-updater.timer
[Unit]
Description=Swiftwave Updater Timer

[Timer]
OnUnitActiveSec=5m
Unit=swiftwave-updater.service

[Install]
WantedBy=timers.target
